# templates/server/compress.timer.epp
<%- | String $on_calendar = 'daily',
      String $randomize_delay_sec = '4h',
| -%>
[Timer]
OnCalendar=<%= $on_calendar %h
<% if $randomize_delay_sec { -%>
RandomizedDelaySec=<%= $randomize_delay_sec %>
<% } -%>

[Install]
WantedBy=timers.target
